# contrib/debian/timbre.service
[Unit]
Description=Timbre self-hosted music server
After=network.target

[Service]
Type=simple
User=timbre
Group=timbre
ExecStart=/opt/timbre/timbre-server
EnvironmentFile=/etc/timbre/timbre.env
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/timbre

[Install]
WantedBy=multi-user.target
